# supervisor.conf
[program:control-page]
directory=/home/pi/work/control-page
command=/bin/bash -E -c ./start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
user=pi
